# conf/node-monitor.service
[Unit]
Description=Node Monitor

[Service]
Environment="HOME=/root"
ExecStart=
Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target
