# conf/node-monitor.service
[Unit]
Description=Node Monitor

[Service]
ExecStart=
Restart=always
RestartSec=120
KillMode=process

[Install]
WantedBy=multi-user.target
